# supervisord.conf
cat > supervisord.conf << 'EOF'
[supervisord]
nodaemon=true
logfile=/data/log/supervisord.log

[program:mongodb]
command=mongod --dbpath /data/db --logpath /data/log/mongod.log --noauth --bind_ip 127.0.0.1
autostart=true
autorestart=true
priority=10
startsecs=5

[program:backend]
command=uvicorn main:app --host 127.0.0.1 --port 8000
directory=/app/backend
autostart=true
autorestart=true
priority=20
startsecs=8
startretries=10
stdout_logfile=/data/log/backend.log
stderr_logfile=/data/log/backend.log

[program:frontend]
command=node server.js
directory=/app/frontend
autostart=true
autorestart=true
priority=30
startsecs=3
stdout_logfile=/data/log/frontend.log
stderr_logfile=/data/log/frontend.log
EOF
