# supervisord.conf
[supervisord]
nodaemon=true
logfile=/data/log/supervisord.log
logfile_maxbytes=10MB
loglevel=info
pidfile=/tmp/supervisord.pid

[program:mongodb]
command=mongod --dbpath /data/db --logpath /data/log/mongod.log --bind_ip 127.0.0.1 --port 27017
autostart=true
autorestart=true
priority=10
startsecs=5
stdout_logfile=/dev/null
stderr_logfile=/dev/null

[program:api]
command=python3 -m uvicorn main:app --host 127.0.0.1 --port 8000
directory=/app/backend
autostart=true
autorestart=true
priority=20
startsecs=5
startretries=15
environment=
    MONGO_URI="mongodb://vsapp:vspassword123@127.0.0.1:27017/vulnsamurai",
    MONGO_DB="vulnsamurai",
    JWT_SECRET="changeme_run_openssl_rand_hex_64_and_paste_here"
stdout_logfile=/data/log/backend.log
stderr_logfile=/data/log/backend.log

[program:frontend]
command=node server.js
directory=/app/frontend
autostart=true
autorestart=true
priority=30
startsecs=3
startretries=10
environment=
    BACKEND_URL="http://localhost:8000",
    PORT="3000",
    SESSION_SECRET="vulnsamurai_session_secret_change_me"
stdout_logfile=/data/log/frontend.log
stderr_logfile=/data/log/frontend.log
